# jenkins.service
[Unit]
Description=Jenkins Automation Server
After=network.target

[Service]
Type=notify
User=jenkins
Group=jenkins

Environment="JENKINS_HOME=/var/lib/jenkins"
Environment="JENKINS_PORT=8080"
Environment="CASC_JENKINS_CONFIG=/var/lib/jenkins/jenkins.yaml"
Environment="JAVA_OPTS=-Djava.awt.headless=true -Xmx2g -Xms512m"

ExecStart=/usr/bin/java $JAVA_OPTS -jar /usr/share/jenkins/jenkins.war --httpPort=$JENKINS_PORT

Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=jenkins

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/var/lib/jenkins /var/log/jenkins

[Install]
WantedBy=multi-user.target
